# nginx-dstat.service
[Unit]
Description=nginx dstat dashboard
After=network.target nginx.service
Requires=nginx.service

[Service]
Type=simple
User=www-data
Group=adm
ExecStart=/usr/local/bin/nginx-dstat \
    -addr :8765 \
    -log /var/log/nginx/access.log \
    -status http://127.0.0.1/nginx_status \
    -interval 1500ms
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nginx-dstat

# Minimal capabilities needed
NoNewPrivileges=true
ProtectSystem=strict
ReadOnlyPaths=/var/log/nginx
PrivateTmp=true

[Install]
WantedBy=multi-user.target
